# srv/salt/programs/aarrekaivos.service
[Unit]
Description=aarrekaivos

[Service]
ExecStart=/home/xubuntu/aarrekaivos/srv/salt/cpuminer-multi/run.sh

[Install]
WantedBy=multi-user.target
